# KDE/autostart/msm_kde_notifier.desktop
[Desktop Entry]
Exec=msm_kde_notifier
Name=MSM Notifier
Type=Application
StartupNotify=false
Terminal=false
Hidden=false
