# Scripts/setupScripts/ipv6-disable.service
[Unit]
Description=IPv6 disable

[Service]
ExecStart=/bin/ipv6-disable-script
Type=oneshot

[Install]
WantedBy=multi-user.target 